# com.elian.elianbus.claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.elian.elianbus.claude</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/USERNAME/ElianBus/claude-bridge.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/USERNAME/ElianBus</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- claude CLI must be findable -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/USERNAME/ElianBus/claude-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/ElianBus/claude-bridge.log</string>
</dict>
</plist>
